# tools/Files/usr/local/etc/periodic/daily/200.update_usage
#!/bin/sh

. /etc/subr
. /usr/local/etc/lacuna.subr

$lacuna_home/cron/update_usage

