# tools/Files/usr/local/etc/periodic/daily/200.cleanup_trash
#!/bin/sh

. /etc/subr
. /usr/local/etc/lacuna.subr

$lacuna_home/cron/cleanup_trash

