# tools/Files/usr/local/etc/periodic/weekly/300.backup
#!/bin/sh

. /etc/subr
. /usr/local/etc/lacuna.subr

$lacuna_home/bin/backup weekly

